# darwin13template/canonical/CASA.app/Contents/Library/LaunchAgents/org.freedesktop.dbus-session.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.freedesktop.dbus-session</string>
	<key>OnDemand</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>/@@PREFIX@@/CASA.app/Contents/MacOS/dbus-daemon</string>
		<string>--nofork</string>
		<string>--config-file=@@PREFIX@@/CASA.app/Contents/Resources/dbus-1/session.conf</string>
	</array>
	<key>ServiceIPC</key>
	<true/>
	<key>Sockets</key>
	<dict>
		<key>unix_domain_listener</key>
		<dict>
			<key>SecureSocketWithKey</key>
			<string>DBUS_LAUNCHD_SESSION_BUS_SOCKET</string>
		</dict>
	</dict>
</dict>
</plist>
